# qiita minecraft/lightsail-graph.service
[Unit]
Description=make server status graph

[Service]
Type=oneshot
ExecStart=/opt/status.py
User=ubuntu

[Install]
WantedBy=multi-user.target
