# qiita minecraft/lightsail-graph.timer
[Unit]
Description=timer lightsail-graph

[Timer]
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
